# aria2/aria2.ini
[program:webui-aria2]
directory=/data0/webui-aria2/
command=python -m SimpleHTTPServer 9999
autostart=true
startsecs=10
autorestart=true
startretries=3
user=root
priority=999
redirect_stderr=false
stdout_logfile=/data0/webui-aria2/logs/aria2-stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_events_enabled=false
stderr_logfile=/data0/webui-aria2/logs/aria2-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_events_enabled=false
stopasgroup=false
killasgroup=false
