# docker/centos-ssh/init.ini
[program:init]
command=/usr/lib/systemd/systemd --system
process_name=%(program_name)s 
autostart=true
