# docker/centos-ssh/sshd.ini
[program:sshd]
command=/usr/sbin/sshd -D
process_name=%(program_name)s 
autostart=true
